# webhook/easybewerbung-webhook.service
[Unit]
Description=EasyBewerbung GitHub Webhook Deployment Service
After=network.target

[Service]
Type=simple
User=alexgiss
Group=alexgiss
WorkingDirectory=/home/alexgiss/EasyBewerbung
Environment=WEBHOOK_SECRET=easybewerbung-deploy-2026-secret
ExecStart=/usr/bin/python3 /home/alexgiss/EasyBewerbung/webhook/deploy.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
